# supervisor/cae-app-scheduler.conf
[program:cae-app-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /home/ubuntu/vision-ca-api/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=ubuntu
numprocs=1
redirect_stderr=true
stdout_logfile=/home/ubuntu/vision-ca-api/storage/logs/scheduler.log
stopwaitsecs=3600
